# supervisor/mysqld.conf
[program:mysqld]
command=/workspace/app/mysql/bin/mysqld_safe
autostart=true
autorestart=true
user=mysql
