# supervisor/mysqld.conf
[program:mysqld]
command=/workspace/app/mysql/bin/mysqld_safe --defaults-file=/workspace/etc/mysql/my.cnf --user=mysql
autostart=true
autorestart=true
user=mysql
