# supervisor/nginx.conf
[program:nginx]
command=/usr/local/sbin/nginx -c /workspace/etc/nginx/nginx.conf
autostart=true
autorestart=true
