# supervisor/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize  --fpm-config /workspace/etc/php/php-fpm.conf
autostart=true
autorestart=true
