# supervisor/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /workspace/etc/redis/6379.conf
autostart=true
autorestart=true
